# config/systemd/iotdb-postgres.service
[Unit]
Description=IoTDB Enhanced PostgreSQL Database
After=network.target
Wants=network-online.target

[Service]
Type=notify
User=postgres
Group=postgres
Environment=PGDATA=/var/lib/postgresql/data
EnvironmentFile=/etc/postgresql/postgresql.conf
ExecStart=/usr/lib/postgresql/15/bin/postgres -D /var/lib/postgresql/data
ExecReload=/bin/kill -HUP $MAINPID
KillMode=mixed
KillSignal=SIGINT
TimeoutSec=infinity
Restart=on-failure
RestartSec=5s

[Install]
WantedBy=multi-user.target
