# config/systemd/iotdb-frontend.service
[Unit]
Description=IoTDB Enhanced Frontend
After=network.target iotdb-backend.service
Wants=iotdb-backend.service

[Service]
Type=forking
User=iotdb
Group=iotdb
WorkingDirectory=/root/iotdb-enhanced/frontend
Environment="NODE_ENV=production"
Environment="PORT=3000"
ExecStart=/usr/bin/pm2 start ecosystem.config.cjs --no-daemon
ExecStop=/usr/bin/pm2 stop ecosystem.config.cjs
Restart=always
RestartSec=10s
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
